# templates/setkey.service
[Unit]
Description=IPSec policy configuration startup task

[Service]
Type=oneshot
ExecStart=/sbin/setkey -f /etc/ipsec-tools.conf

[Install]
WantedBy=multi-user.target
